# roles/openvpn/files/etc/systemd/system/openvpn_exporter.service
[Unit]
Description=OpenVPN Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=openvpn
ExecStart=/usr/local/bin/openvpn_exporter -openvpn.status_paths /var/log/openvpn/openvpn-status.log -web.listen-address :29180

[Install]
WantedBy=default.target
